# etc/cron.d/duckdns
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
*/5 * * * * root duckdns >/dev/null 2>&1
